# virtfs/cvsccomc/etc/profile.d/bash_timestamps.sh

# Enable timestamps in bash history
export HISTTIMEFORMAT="%F %T "

